# systemd/backup.service
[Unit]
Description=backup

[Service]
Type=oneshot
ExecStart=/usr/bin/rsync /home/franz/work -v -a --exclude "gocode/pkg" /mnt/backup
