# systemd/backup.service
[Unit]
Description=backup

[Service]
Type=oneshot
ExecStart=/usr/bin/rsync /home/franz/work -a /mnt/backup
